# src/postgresql.conf
[program:postgresql]
command=postgres -D /data/postgresql -c unix_socket_directory=/var/run/postgresql.sock
autostart=true
autorestart=true
startretries=5
