# src/pxls.conf
[program:pxls]
command="java -jar pxls*.jar"
autostart=true
autorestart=true
startretries=5
